# export/lix_export_worker.conf
[program: lix_export_worker]
command=/opt/conda_envs/bluesky_workers/bin/python /home/xf16id/lix-workers/export/lix_export_worker.py
user=xf16id